# atlas/services/unbound.service.d/override.conf
[Service]
OOMScoreAdjust=-500
Restart=on-failure
RestartSec=10s