# hermes/services/mosquitto.service.d/01-run_after_unbound.conf
[Unit]
Requires=unbound.service
After=network.target unbound.service

