# lilith/services/blocklist-update.service
[Unit]
Description=Update the blocklist
After=network.target

[Service]
Type=simple
ExecStart=/root/tpl/blocklist-update.sh

