# atlas/services/ntpd.service.d/override.conf
[Unit]
StartLimitInterval=2min
StartLimitBurst=10

[Service]
ExecStartPre=/usr/bin/ntpdate -b 129.6.15.28
ExecStopPost=/usr/bin/hwclock -w
RestartSec=10
