# atlas/services/ntpd.service.d/override.conf
[Service]
ExecStartPre=/usr/bin/ntpdate -b 129.6.15.28
ExecStopPost=/usr/bin/hwclock -w