# atlas/services/thttpd.service.d/override.conf
[Service]
OOMScoreAdjust=500
Restart=on-failure
RestartSec=60s